# motd_speedtest.timer
[Unit]
Description=Generate the Message Of The Day

[Timer]
OnBootSec=6min
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
